# src/boiler/etc/cron.d/heating
# M H DoM M DoW Cmd
*/5 * * 10,11,12,1,2,3 * root /opt/heating/montemp

